# docker/centos/curl.conf
[program:curl]
command = curl http://169.254.169.254/latest/meta-data/public-keys/0/openssh-key -o /home/centos/.ssh/authorized_keys
priority = 2
startretries = 60
user = centos
