# docker/centos/sshd-keygen.conf
[program:sshd-keygen]
command = /usr/sbin/sshd-keygen
priority = 0
startsecs = 0
